# roles/motioneye/templates/motioneye.service.j2
[Unit]
Description=motionEye Server

[Service]
ExecStart=/usr/local/bin/meyectl startserver -c {{ motioneye_config_dir }}/motioneye.conf
Restart=on-abort

[Install]
WantedBy=multi-user.target
