# roles/glances/files/glances.service
# GENERATED BY ANSIBLE
[Unit]
Description=Glances webserver
Documentation=https://glances.readthedocs.io/en/latest/
After=network.target

[Service]
ExecStart=/usr/local/bin/glances -w
Restart=on-abort
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
